# klax.service
[Unit]
Description=klax — AI messaging bridge
After=network.target
StartLimitBurst=3
StartLimitIntervalSec=60

[Service]
Type=simple
ExecStart=%h/.local/bin/klax start --foreground
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
